# supervisor/api.ini
[program:laravel_new_api_event_listener]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/${VHOST}/current/artisan event-worker
user=nginx
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/laravel_new_api_event_listener.err.log
stdout_logfile=/var/log/supervisor/laravel_new_api_event_listener.out.log

[program:laravel_new_api_communication_listener]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/${VHOST}/current/artisan communication-worker
user=nginx
numprocs=1
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/laravel_new_api_communication_listener.err.log
stdout_logfile=/var/log/supervisor/laravel_new_api_communication_listener.out.log